# modules/tectonic/resources/tectonic.path
[Unit]
Description=Trigger for tectonic.service
[Path]
PathExists=/opt/tectonic/tectonic-wrapper.sh
Unit=tectonic.service
[Install]
WantedBy=multi-user.target
