# data/data/bootstrap/aro/files/etc/NetworkManager/dispatcher.d/30-eth0-mtu-3900.template
#!/bin/bash

if [ "$1" == "eth0" ] && [ "$2" == "up" ]; then
    ip link set $1 mtu 3900
fi
