# data/data/bootstrap/systemd/common/units/systemd-journal-gatewayd.service.d/certs.conf
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-journal-gatewayd \
  --key=/opt/openshift/tls/journal-gatewayd.key \
  --cert=/opt/openshift/tls/journal-gatewayd.crt \
  --trust=/opt/openshift/tls/root-ca.crt
