# data/data/bootstrap/files/etc/systemd/system/node-image-pivot.target
[Unit]
Description=Node Image Pivot Target
Requires=basic.target

Requires=node-image-pivot.service
