# modules/bootkube-ut2/resources/bootkube.path
[Unit]
Description=Trigger for bootkube.service
[Path]
PathExists=/opt/tectonic/bootkube.sh
Unit=bootkube.service
[Install]
WantedBy=multi-user.target
